# systemd/liverootsafety.target
[Unit]
Description=unmount /boot if running rootFS in ram
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes

